# rootfs/etc/services.d/webui/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "URnetwork Web UI stopped"
